# chuchu_lobby@.service.in
[Unit]
Description=%i lobby server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/chuchu_lobby_server SYSCONFDIR/%i.cfg
StandardOutput=append:LOCALSTATEDIR/log/%i-lobby.log

[Install]
WantedBy=multi-user.target
