# chuchu_login@.service.in
[Unit]
Description=%i login server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/chuchu_login_server SYSCONFDIR/%i.cfg
StandardOutput=append:LOCALSTATEDIR/log/%i-login.log

[Install]
WantedBy=multi-user.target
